# templates/ntp-stats.cron.j2
# {{ ansible_managed }}
*/10 * * * * root sleep 23 && nice -10 {{ helper_path }} ntp-stats > /dev/null 2>&1
